# services/selenium/etc/service/selenium/run
#!/bin/bash
eval /opt/bin/entry_point.sh
