# services/execution_engine/etc/service/cromwell/run
#!/bin/bash
eval java -jar /opt/cromwell/target/scala-2.11/cromwell-0.10.jar server
